# connector/systemd/zecurity-connector-update.service
# zecurity-connector-update.service — Oneshot auto-update check
#
# Triggered daily by zecurity-connector-update.timer.
# Runs as root because it needs to replace /usr/local/bin/zecurity-connector
# and call systemctl to restart the main service.
#
# The main zecurity-connector.service runs as the unprivileged 'zecurity' user
# and cannot perform these operations itself.

[Unit]
Description=ZECURITY Connector Auto-Update Check
Documentation=https://github.com/vairabarath/zecurity
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
EnvironmentFile=-/etc/zecurity/connector.conf
Environment=AUTO_UPDATE_ENABLED=true
ExecStart=/usr/local/bin/zecurity-connector --check-update
